# ds-sensors-mqtt/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock ; the path to the socket file

[inet_http_server]
port=0.0.0.0:9001 ; ip_address:port specifier, *:port for all iface

[supervisord]
logfile=%(here)s/supervisord.log ; main log file; default $CWD/supervisord.log
loglevel=info                   ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid     ; supervisord pidfile; default supervisord.pid
nodaemon=true                    ; start in foreground if true; default false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
;serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
serverurl=http://0.0.0.0:9001          ; use an http:// url to specify an inet socket

[program:ds-sensors-mqtt]
autorestart=true
autostart=true
command=make start
redirect_stderr=true
stdout_logfile=%(here)s/program.log
